# docker/slash/etc/supervisord.conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
numprocs=1
autostart=true
startsecs=1
startretries=3
stdout_logfile=/var/log/nginx/evaneos_output.log
stderr_logfile=/var/log/nginx/evaneos_error.log