# docker/slash/etc/supervisord.conf.d/nodejs.conf
[program:nodejs]
command=node /welean3-tada/node/index.js
autostart=true
autorestart=true
user=root
priority=100
stdout_logfile=/var/log/stdout_nodejs.log
stderr_logfile=/var/log/stderr_nodejs.log