# docker/slash/etc/supervisord.conf.d/ssh.conf
[program:sshd]
command=/usr/sbin/sshd -D
numprocs=1
autostart=true
startsecs=1
startretries=3
stdout_logfile=/var/log/sshd_output.log
stderr_logfile=/var/log/sshd_error.log